# contrib/ssh-mcp-daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  LaunchAgent for the ssh-mcp daemon.

  Replace /ABSOLUTE/PATH/TO/ssh-mcp with the absolute path of your installed
  binary, copy this file to ~/Library/LaunchAgents/, then run:

    launchctl load ~/Library/LaunchAgents/ssh-mcp-daemon.plist

  If the daemon cannot reach your SSH agent, add an <EnvironmentVariables>
  dict that sets SSH_AUTH_SOCK to its socket path.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ssh-mcp-daemon</string>
    <key>ProgramArguments</key>
    <array>
        <string>/ABSOLUTE/PATH/TO/ssh-mcp</string>
        <string>daemon</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/tmp/ssh-mcp-daemon.log</string>
</dict>
</plist>
